# vesyl-print-agent.service
[Unit]
Description=VESYL Print — cloud agent (heartbeat / pairing / OTA)
After=network-online.target cups.service
Wants=network-online.target

[Service]
Type=simple
User=vesyl
WorkingDirectory=/opt/vesyl-print/current
ExecStart=/usr/bin/python3 /opt/vesyl-print/current/agent.py
Restart=on-failure
RestartSec=5
# Do not put secrets in Environment= lines.
Environment=PYTHONUNBUFFERED=1
Environment=VESYL_PRINT_INSTALL_ROOT=/opt/vesyl-print

[Install]
WantedBy=multi-user.target
