# vesyl-print-agent.service
[Unit]
Description=VESYL Print — cloud agent (heartbeat / pairing)
After=network-online.target cups.service
Wants=network-online.target

[Service]
Type=simple
User=vesyl
WorkingDirectory=/home/vesyl/vesyl-print
ExecStart=/usr/bin/python3 /home/vesyl/vesyl-print/agent.py
Restart=on-failure
RestartSec=5
# Do not put secrets in Environment= lines.
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
