# vesyl-print-display.service
[Unit]
Description=VESYL Print — LCD system info display
# Wait until Plymouth has quit before painting. Plymouth owns the LCD during
# boot (showing the VESYL splash) and repaints it periodically; if the app
# starts while Plymouth is still up, the two fight over fb1 and the screen
# flickers between the splash and the live view. Ordering after plymouth-quit
# gives a clean hand-off: splash through boot, then the app takes over.
After=plymouth-quit-wait.service

[Service]
Type=simple
User=vesyl
WorkingDirectory=/opt/vesyl-print/current
ExecStart=/usr/bin/python3 /opt/vesyl-print/current/main.py
ExecStopPost=/usr/bin/python3 /opt/vesyl-print/current/main.py --offline
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
